# node/casperlabs-node.service
[Unit]
Description=CasperLabs Node
After=network.target

[Service]
ExecStart=/usr/bin/casperlabs-node run --data-dir=/var/lib/casperlabs-node
User=node

[Install]
WantedBy=multi-user.target
